# 90-keyboard-backlight.rules
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="smc::kbd_backlight", RUN+="/bin/chmod 0666 /sys/class/leds/%k/brightness"